# deploy/systemctl/rocketmq-broker.service
[Unit]
Description=Apache RocketMQ Broker (Single Mode)
After=network.target docker.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/rocketmq-4.9.8
ExecStart=/opt/rocketmq-4.9.8/bin/mqbroker -c /opt/rocketmq-4.9.8/conf/single-broker.conf
Restart=on-failure
RestartSec=10
StandardOutput=append:/data/logs/broker.log
StandardError=append:/data/logs/broker.log

[Install]
WantedBy=multi-user.target